# appliance/gmktec/docker.service.d/10-mbfd-storage.conf
[Unit]
# Docker's data-root is /mnt/mbfd-storage/docker-data. Never let the daemon
# initialize against the hidden root-disk mountpoint during a slow disk mount.
RequiresMountsFor=/mnt/mbfd-storage
Wants=network-online.target
After=network-online.target
